# ceph-log-rotate.service
[Unit]
Description=Ceph log rotation

[Service]
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/ceph
Type=oneshot
