# crontab
* * * * * systemctl start --no-block ceph-log-rotate.service
